# docker/supervisor/layered/development/watch.conf
[program:watch]
directory=/gjsx
command=/bin/sh -c "yarn watch"
stderr_logfile = /gjsx/docker/log/dev/watch_error.log
stdout_logfile = /gjsx/docker/log/dev/watch.log
autorestart=false


