# docker/supervisor/prod/broadway.conf
[program:broadwayd]
command=/bin/sh -c "nohup gtk4-broadwayd $BROADWAY_DISPLAY"
stderr_logfile = /var/log/broadwayd.log
stdout_logfile = /var/log/broadwayd.log
autorestart=false
priority=100
